# conf/supervisor.vilogged.conf
[program:rapidpro]
command=/usr/bin/uwsgi --ini /code/conf/uwsgi.ini
directory=/code/
autostart=true
autorestart=true
stdout_logfile=/var/log/vilogged/vilogged.log
redirect_stderr=True