# conf/supervisor.vilogged-cron.conf
[program:vcron]
command=python  /code/vilogged/task.py
directory=/code/
autostart=true
autorestart=true
stdout_logfile=/var/log/vilogged/vilogged.log
redirect_stderr=True